# backend/immobrowse.service
[Unit]
Description=ImmoBrowse

[Service]
User=immobrowse
Group=immobrowse
ExecStart=/usr/local/bin/immobrowsed

[Install]
WantedBy=wsgi-services.target
